# engram-server.service
# Engram v3 MCP Server -- systemd user service
#
# Runs engram in SSE mode so remote editors can connect over the network.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp engram-server.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now engram-server
#
# Logs:
#   journalctl --user -u engram-server -f

[Unit]
Description=Engram v3 MCP Memory Server (SSE + Postgres)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/engram
ExecStart=%h/engram/.venv/bin/python -m engram --transport sse --port 8788
Restart=on-failure
RestartSec=5

Environment=PYTHONPATH=%h/engram/src

# Load DATABASE_URL, OPENAI_API_KEY, etc. from a secrets file (recommended)
EnvironmentFile=-%h/.config/engram/env

[Install]
WantedBy=default.target
